# etc/xdg/autostart/myapp.desktop
[Desktop Entry]
Exec=chromium-browser --start-fullscreen --start-maximized --bwsi http://domain/RaspberryPints/
