# etc/xdg/autostart/myapp.desktop
[Desktop Entry]
Exec=chromium-browser --start-fullscreen --start-maximized --bwsi http://10.0.0.185/RaspberryPintsLiving/
