# etc/xdg/autostart/myapp.desktop
[Desktop Entry]
Exec=chromium-browser --start-fullscreen http://10.0.0.185/RaspberryPintsLiving/
