# etc/xdg/autostart/myapp.desktop
[Desktop Entry]
Exec=chromium-browser --start-fullscreen --start-maximized --incognito --kiosk http://domain/RaspberryPints/